# Config/crontab
@daily zsh -c 'zinit update --parallel $(nproc) && zinit self-update'
